# systemd/rider-voice.service
[Unit]
Description=Rider-Pi Voice Service
After=network-online.target sound.target

[Service]
Type=simple
WorkingDirectory=/home/pi/robot
ExecStart=/usr/bin/python3 -u -m apps.voice.service
Restart=on-failure
Environment=PYTHONUNBUFFERED=1
Environment=VOICE_BIND=127.0.0.1:8093
Environment=VOICE_ASR_BACKEND=stub
TimeoutStopSec=3s

[Install]
WantedBy=multi-user.target
