# systemd/rider-edge-preview.service
[Unit]
Description=Rider-Pi EDGE camera preview (Canny -> snapshots/proc.jpg)
After=network-online.target
Wants=network-online.target rider-api.service
Conflicts=rider-ssd-preview.service rider-cam-preview.service
StartLimitIntervalSec=30
StartLimitBurst=10

[Service]
Environment=SNAP_EVERY_MS=500
Environment=FRAME_W=640
Environment=FRAME_H=480
Type=simple
Environment=SNAP_EVERY_MS=500
Environment=FRAME_W=640
Environment=FRAME_H=480
User=pi
Environment=SNAP_EVERY_MS=500
Environment=FRAME_W=640
Environment=FRAME_H=480
Group=pi
Environment=SNAP_EVERY_MS=500
Environment=FRAME_W=640
Environment=FRAME_H=480
WorkingDirectory=/home/pi/robot
Environment=SNAP_EVERY_MS=500
Environment=FRAME_W=640
Environment=FRAME_H=480
EnvironmentFile=-/etc/default/rider
Environment=SNAP_EVERY_MS=500
Environment=FRAME_W=640
Environment=FRAME_H=480
EnvironmentFile=-/etc/default/rider-edge-preview
Environment=SNAP_EVERY_MS=500
Environment=FRAME_W=640
Environment=FRAME_H=480

Environment=SNAP_EVERY_MS=500
Environment=FRAME_W=640
Environment=FRAME_H=480
# Domyślne parametry (możesz nadpisać w /etc/default/rider-edge-preview)
Environment=SNAP_EVERY_MS=500
Environment=FRAME_W=640
Environment=FRAME_H=480
Environment=SNAP_DIR=/home/pi/robot/snapshots
Environment=SNAP_EVERY_MS=500
Environment=FRAME_W=640
Environment=FRAME_H=480
Environment=SNAP_EVERY_MS=500
Environment=FRAME_W=640
Environment=FRAME_H=480
Environment=EDGE_LOW=60
Environment=SNAP_EVERY_MS=500
Environment=FRAME_W=640
Environment=FRAME_H=480
Environment=EDGE_HIGH=120
Environment=SNAP_EVERY_MS=500
Environment=FRAME_W=640
Environment=FRAME_H=480
Environment=PREVIEW_ROT=0
Environment=SNAP_EVERY_MS=500
Environment=FRAME_W=640
Environment=FRAME_H=480
Environment=PREVIEW_FLIP_H=0
Environment=SNAP_EVERY_MS=500
Environment=FRAME_W=640
Environment=FRAME_H=480
Environment=FRAME_W=640
Environment=SNAP_EVERY_MS=500
Environment=FRAME_H=480
Environment=FRAME_H=480
Environment=SNAP_EVERY_MS=500
Environment=FRAME_W=640
Environment=LAST_FRAME=/home/pi/robot/data/last_frame.jpg
Environment=SNAP_EVERY_MS=500
Environment=FRAME_W=640
Environment=FRAME_H=480

Environment=SNAP_EVERY_MS=500
Environment=FRAME_W=640
Environment=FRAME_H=480
# Fallback (żeby /camera/proc nie sypał 404 na starcie) i katalogi
Environment=SNAP_EVERY_MS=500
Environment=FRAME_W=640
Environment=FRAME_H=480
ExecStartPre=/bin/mkdir -p /home/pi/robot/snapshots /home/pi/robot/data
Environment=SNAP_EVERY_MS=500
Environment=FRAME_W=640
Environment=FRAME_H=480
ExecStartPre=/bin/ln -sf /home/pi/robot/snapshots/raw.jpg /home/pi/robot/snapshots/proc.jpg
Environment=SNAP_EVERY_MS=500
Environment=FRAME_W=640
Environment=FRAME_H=480

Environment=SNAP_EVERY_MS=500
Environment=FRAME_W=640
Environment=FRAME_H=480
# Lock na kamerę + start
Environment=SNAP_EVERY_MS=500
Environment=FRAME_W=640
Environment=FRAME_H=480
ExecStart=/usr/bin/flock -n /tmp/camera.lock /usr/bin/python3 -u apps/vision/edge_preview.py
Environment=SNAP_EVERY_MS=500
Environment=FRAME_W=640
Environment=FRAME_H=480
Environment=SNAP_EVERY_MS=500
Environment=FRAME_W=640
Environment=FRAME_H=480

Environment=SNAP_EVERY_MS=500
Environment=FRAME_W=640
Environment=FRAME_H=480
Restart=always
Environment=SNAP_EVERY_MS=500
Environment=FRAME_W=640
Environment=FRAME_H=480
RestartSec=0.5
Environment=SNAP_EVERY_MS=500
Environment=FRAME_W=640
Environment=FRAME_H=480
TimeoutStartSec=10
Environment=SNAP_EVERY_MS=500
Environment=FRAME_W=640
Environment=FRAME_H=480
TimeoutStopSec=5
Environment=SNAP_EVERY_MS=500
Environment=FRAME_W=640
Environment=FRAME_H=480
Environment=SNAP_EVERY_MS=500
Environment=FRAME_W=640
Environment=FRAME_H=480
KillSignal=SIGINT
Environment=SNAP_EVERY_MS=500
Environment=FRAME_W=640
Environment=FRAME_H=480
Nice=-5
Environment=SNAP_EVERY_MS=500
Environment=FRAME_W=640
Environment=FRAME_H=480
StandardOutput=journal
Environment=SNAP_EVERY_MS=500
Environment=FRAME_W=640
Environment=FRAME_H=480
StandardError=journal
Environment=SNAP_EVERY_MS=500
Environment=FRAME_W=640
Environment=FRAME_H=480

Environment=SNAP_EVERY_MS=500
Environment=FRAME_W=640
Environment=FRAME_H=480
[Install]
WantedBy=multi-user.target
