# systemd/rider-ui-manager.service
[Unit]
Description=Rider-Pi UI Manager (pause/dim/standby)
After=rider-broker.service rider-motion.service rider-camera.service
Wants=rider-broker.service

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/robot
Environment=BUS_SUB_ADDR=tcp://127.0.0.1:5556
Environment=BUS_PUB_ADDR=tcp://127.0.0.1:5555
Environment=MOTION_STATE_TOPIC=motion.state
Environment=VISION_TOPIC=vision.state
Environment=UI_CTRL_TOPIC=ui.control
Environment=UI_CHAT_MODE=1
Environment=UI_INACTIVITY_DIM_SEC=30
Environment=UI_INACTIVITY_OFF_SEC=120
Environment=UI_DIM_MODE=xgo
Environment=UI_XGO_DIM=10
Environment=UI_XGO_BRIGHT=80
Environment=UI_XGO_BLACK_DIM=0
Environment=UI_AUDIO_DIM_PCT=20
Environment=UI_AUDIO_OFF_MUTE=1
ExecStart=/usr/bin/python3 -u /home/pi/robot/apps/ui/manager.py
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target