# systemd/rider-camera-preview.service
[Unit]
Description=Rider-Pi camera preview to 2" SPI LCD (takeover)
After=local-fs.target
StartLimitIntervalSec=0

[Service]
Type=simple
User=root
WorkingDirectory=/home/pi/robot
Environment=PYTHONUNBUFFERED=1
Environment=VISION_HUMAN=1
Environment=VISION_FACE_EVERY=5
Environment=PREVIEW_ROT=270
ExecStart=/usr/bin/python3 -m apps.camera
ExecStopPost=/usr/bin/python3 /home/pi/robot/ops/lcdctl.py off
Restart=on-failure
RestartSec=1

[Install]
WantedBy=multi-user.target
