# systemd/rider-web-bridge.service
[Unit]
Description=Rider-Pi HTTP→ZMQ Motion Web Bridge
After=network-online.target rider-broker.service
Wants=network-online.target rider-broker.service

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/robot
# Produkcyjny serwer WSGI (zamiast Flask dev servera)
ExecStart=/usr/bin/python3 -u -m services.web_motion_bridge
RestartSec=0.5
StandardOutput=journal
StandardError=journal

# HTTP
Environment=WEB_T_DEFAULT=0.10
Environment=WEB_V_DEFAULT=0.08
Environment=WEB_HOST=0.0.0.0
Environment=WEB_PORT=8081
# Rate-limit po stronie HTTP mostka
Environment=WEB_MIN_GAP_MS=90
# PUB: łączymy się do brokera (mostek jako klient PUB)
Environment=WEB_BIND_PUB=0

# BUS
Environment=BUS_PUB_PORT=5555
Environment=BUS_SUB_PORT=5556
Environment=BUS_RCVHWM=100
Environment=BUS_CONFLATE=0

[Install]
WantedBy=multi-user.target

