# systemd/rider-api.service
[Unit]
Description=Rider-Pi Status API (Flask)
After=network.target rider-broker.service
Wants=rider-broker.service

[Service]
User=pi
Group=pi
WorkingDirectory=/home/pi/robot
Environment=PYTHONUNBUFFERED=1
Environment=API_HOST=0.0.0.0
Environment=API_PORT=8080
ExecStart=/usr/bin/python3 -u /home/pi/robot/services/status_api.py
Restart=always
RestartSec=1

[Install]
WantedBy=multi-user.target
